# deploy/minirouter.service
# DEPRECATED: This systemd service is for bare-metal / non-Docker deployments.
# MiniRouter now uses Docker as the primary deployment method.
# See docker-compose.yml for the recommended setup.

[Unit]
Description=MiniRouter — Smart LLM Routing Gateway (deprecated: use Docker)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=minirouter
Group=minirouter
WorkingDirectory=/opt/minirouter

# Run from the built dist
ExecStart=/usr/bin/node /opt/minirouter/dist/serve.js

# Graceful shutdown
KillSignal=SIGTERM
TimeoutStopSec=30

# Logging — journald, keep last 100M
StandardOutput=journal
StandardError=journal

# Restart policy
Restart=always
RestartSec=5

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false

[Install]
WantedBy=multi-user.target
